# infra/systemd/willbuy-visitor-worker.service
[Unit]
# willbuy visitor worker — polls DB for studies in `visiting` status,
# calls LLM, and transitions studies visiting→aggregating.
#
# Depends on the API being up (shares app.env and the same DB).
# Entry point: apps/visitor-worker/src/bin.ts (the canonical binary;
# the old scripts/run-visitor-worker-dogfood.ts referenced from a working
# copy was untracked in 517 — referencing it would crash-loop the unit).

Description=willbuy visitor worker
After=network.target willbuy-api.service willbuy-capture-broker.service

[Service]
Type=simple
User=willbuy
WorkingDirectory=/srv/willbuy
EnvironmentFile=/etc/willbuy/app.env
ExecStart=/home/willbuy/.bun/bin/bun run apps/visitor-worker/src/bin.ts
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
