# infra/systemd/willbuy-capture-broker.service
[Unit]
# willbuy capture-broker — spec §5.13.
#
# Listens on /run/willbuy/broker.sock for ONE typed request per capture.
# Socket inode mode is 0660 (owner root, group willbuy). The mode is
# enforced by chmodSync(socketPath, 0o660) in startBroker() immediately
# after server.listen() — NOT by this unit. UMask=0007 would produce 0770
# without the explicit chmod.
#
# This unit is a TEMPLATE — production reads its env from
# /etc/willbuy/broker.env (provisioned out-of-band via 1Password +
# scripts/push-secrets.sh; the env file is NEVER in git).
#
# N2 note: worker-side bind-mount (READ-ONLY mount of broker.sock into
# the capture container at a fixed in-container path) is tracked in a
# follow-up issue and lands with the capture-worker integration.

Description=willbuy capture broker (spec §5.13)
After=network.target

[Service]
Type=simple
User=root
Group=willbuy
UMask=0007
RuntimeDirectory=willbuy
RuntimeDirectoryMode=0750
EnvironmentFile=-/etc/willbuy/broker.env
ExecStart=/usr/bin/node /opt/willbuy/capture-broker/dist/bin.js
Restart=on-failure
RestartSec=2

# Hardening — defense in depth around the host-side process. The
# container hardening lives in the per-capture container manifest
# (§2 #2 + apps/capture-worker/Dockerfile).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ReadWritePaths=/run/willbuy

[Install]
WantedBy=multi-user.target
