# infra/systemd/willbuy-capture-broker.service
[Unit]
# willbuy capture-broker — spec §5.13.
# Listens on /run/willbuy/broker.sock for ONE typed request per capture.
# Required dependency for willbuy-capture-worker — without this service
# every visit fails with broker.send_failed.

Description=willbuy capture broker (spec §5.13)
After=network.target willbuy-api.service

[Service]
Type=simple
User=root
Group=willbuy
UMask=0007
RuntimeDirectory=willbuy
RuntimeDirectoryMode=0750
WorkingDirectory=/srv/willbuy
EnvironmentFile=/etc/willbuy/app.env
EnvironmentFile=-/etc/willbuy/broker.env
ExecStart=/home/willbuy/.bun/bin/bun run apps/capture-broker/src/bin.ts
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
